# collect_advice.service
[Unit]
Description=Collect random advice from API

[Service]
Type=oneshot
ExecStart=%h/cs1500advices/collect_advice.sh
